# deploy/supervisor/djpmp_prd.ini
[program:djpmp_prd]
command=/data/prd/djpmp/env/bin/gunicorn -c gunicorn.conf.py -p gunicorn-djpmp.pid  config.wsgi
directory=/data/prd/djpmp/src
user=pyer
autostart=true
autorestart=true
startretries=5
stdout_logfile=/data/prd/djpmp/logs/supervisor_stdout.log
stderr_logfile=/data/prd/djpmp/logs/supervisor_stderr.log
environment=DEBUG="False"